# scripts/hf-auto-mount.service
[Unit]
Description=Auto-mount HuggingFace repos for Phy-lab
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/weed/psedulab/curation_tools/scripts/hf_auto_mount.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
